# daemon/dev.touchbridge.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.touchbridge.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/touchbridged</string>
        <string>serve</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>Library/Logs/TouchBridge/daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>Library/Logs/TouchBridge/daemon.stderr.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>NSBluetoothAlwaysUsageDescription</key>
    <string>TouchBridge daemon uses Bluetooth to communicate with paired iOS devices.</string>
</dict>
</plist>
